# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mqtt2influxdb]
command=python3 -u /app/mqtt_bridge.py
autostart=true
autorestart=true
